# dataset/tools/repo-ipfs-wiki/wai.cnwiki.service
[Unit]
Description=Wator Wai Text Ipfs cn
[Service]
Type=simple
ExecStart=/ceph/storage3/ipfs/wai.native.wator/dataset/tools/repo-ipfs-wiki/cnwiki.sh
ExecStop=/bin/kill ${MAINPID}
Restart=always
RestartSec=10s

[Install]
WantedBy=default.target
