# dataset/tools/ipfs-wikidump/wai.zhwiki.service
[Unit]
Description=Wator Wai Text zh

[Service]
Type=simple
ExecStart=/watorvapor/wai.storage/wai.native.wator/ipfs-wikidump/zhwiki.sh
ExecStop=/bin/kill ${MAINPID}
Restart=always
RestartSec=10s

[Install]
WantedBy=default.target
