# ipfs-wikidump/wai.jawiki.service
[Unit]
Description=Wator wai jawiki to IPFS Storage 

[Service]
Type=simple
ExecStart=/watorvapor/wai.storage/wai.native.wator/ipfs-wikidump/jawiki.sh
ExecStop=/bin/kill ${MAINPID}
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target
